# nexarena-scheduler.service
# =============================================================================
# Nexarena - Service systemd pour le Symfony Scheduler
# =============================================================================
# Ce fichier gere le worker Messenger qui execute les taches planifiees.
#
# Installation :
#   1. Copier ce fichier dans /etc/systemd/system/
#      sudo cp nexarena-scheduler.service /etc/systemd/system/
#
#   2. Adapter les chemins si necessaire (WorkingDirectory, ExecStart)
#
#   3. Recharger systemd :
#      sudo systemctl daemon-reload
#
#   4. Activer le demarrage automatique :
#      sudo systemctl enable nexarena-scheduler
#
#   5. Demarrer le service :
#      sudo systemctl start nexarena-scheduler
#
#   6. Verifier le statut :
#      sudo systemctl status nexarena-scheduler
#
#   7. Voir les logs :
#      sudo journalctl -u nexarena-scheduler -f
#
# =============================================================================

[Unit]
Description=Nexarena Symfony Scheduler Worker
After=network.target mysql.service mariadb.service
Wants=mysql.service mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/nexarena
ExecStart=/usr/bin/php bin/console messenger:consume scheduler_default --time-limit=3600 --memory-limit=128M
Restart=always
RestartSec=5

# Securite
NoNewPrivileges=true
PrivateTmp=true

# Environnement
Environment=APP_ENV=prod

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexarena-scheduler

[Install]
WantedBy=multi-user.target
